# mercuriod/contrib/mercuriod.service
[Unit]
Description=Mercurio MQTT Broker
Documentation=https://github.com/frenzox/mercurio
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/mercuriod --config /etc/mercurio/config.toml
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Allow binding to privileged ports if needed
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mercuriod

[Install]
WantedBy=multi-user.target
